# ansible/roles/litecoin/templates/litecoin.service.j2
[Unit]
Description=Litecoin daemon serivce
After=network.target

[Service]
Type=simple
User=ubuntu
ExecStart=/home/ubuntu/.nvm/versions/node/v4.8.7/bin/litecoind
[Install]
WantedBy=multi-user.target